# deploy/systemd/kuro-research-center.service.example
[Unit]
Description=Kuro Research Center
After=network.target

[Service]
WorkingDirectory=/home/kuro/projects/kuro-research-center
EnvironmentFile=/home/kuro/projects/kuro-research-center/.env
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port 8443 --ssl-keyfile certs/key.pem --ssl-certfile certs/cert.pem
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
